# dotfilesd/service/dotfilesd.service.template
[Unit]
Description=dotfiles runtime daemon (Connect RPC + MCP)
Documentation=https://github.com/gpakosz/.tmux
After=network.target

[Service]
Type=simple
ExecStart={{BINARY_DIR}}/dotfilesd
Restart=on-failure
RestartSec=5
Environment=DOTFILESD_PORT=9105
Environment=DOTFILESD_MCP_PORT=9106

[Install]
WantedBy=default.target
