# dotfilesd/service/dotfilesd.service.template
[Unit]
Description=dotfiles runtime daemon (Connect RPC)
After=network.target

[Service]
Type=simple
ExecStart={{BINARY_DIR}}/dotfilesd
Restart=on-failure
RestartSec=5
Environment=DOTFILESD_PORT=9105

[Install]
WantedBy=default.target
